# VagrantApp/uwsgi.conf
description "uwsgi tiny instance"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec uwsgi --socket 10.90.90.95:4242 --chdir /vagrant/ -w wsgi

